# services/system/boot-notifier/boot-notifier.service
[Unit]
Description=Boot Notification Service
Documentation=https://github.com/YourRepo/web1
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
ExecStart=/usr/local/bin/boot-notifier.sh
StandardOutput=journal
StandardError=journal
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
